# systemd/spotpi.service
[Unit]
Description=SpotPi Web UI
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=spotpi
Group=spotpi
SupplementaryGroups=audio
Environment=PCS_CONFIG=/etc/spotpi/config.toml
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/spotpi/venv/bin/spotpi
Restart=always
RestartSec=5
TimeoutStartSec=30
TimeoutStopSec=20
UMask=0077
PrivateTmp=true
ProtectHome=true
ProtectSystem=full
LockPersonality=true
RestrictRealtime=true
ReadWritePaths=/etc/spotpi /var/cache/spotpi /var/lib/spotpi

[Install]
WantedBy=multi-user.target
